# systemd/waiverdb.socket
[Unit]
Description=WaiverDB Socket

[Socket]
ListenStream=/run/waiverdb/socket

[Install]
WantedBy=sockets.target
